# arch/systemd/cadence.service
[Unit]
Description=Cadence

[Service]
ExecStart=cadence-session-start --system-start

[Install]
WantedBy=default.target
