# arch/systemd/discord.service
[Unit]
Description=Discord

[Service]
ExecStart=apulse discord

[Install]
WantedBy=default.target
