# arch/systemd/jack-matchmaker.service
[Unit]
Description=Jack MatchMaker

[Service]
ExecStart=jack-matchmaker -p ~/dotfiles/arch/audio/jack-matchmaker

[Install]
WantedBy=default.target
