# arch/systemd/synergy.service
[Unit]
Description=Synergy

[Service]
ExecStart=synergys -f

[Install]
WantedBy=default.target

